# br2/board/pc-x86_64-uefi/profile.d/dosemu2.sh
dosemu --Fdrive_c /mnt/drive_c -term -T -I 'ignore_djgpp_null_derefs off'
